# systemd/promoradar-relatorio.service
[Unit]
Description=PromoRadar Relatorio Diario
After=network.target

[Service]
Type=simple
User=__PROMORADAR_USER__
WorkingDirectory=__PROMORADAR_WORKDIR__
EnvironmentFile=/etc/default/promoradar
ExecStart=/usr/bin/env bash -c 'cd __PROMORADAR_WORKDIR__ && PYTHONPATH=. python3 scripts/agendar_relatorio_diario.py'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
